# service.conf
[program:Service]
directory=/root/
command=/root/admin.sh
user=user_to_run_program
autostart=true
autorestart=false
stderr_logfile=/path/to/error.log
stdout_logfile=/path/to/console.out.log
